# contrib/init/zebnad.conf
description "Zebna Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env ZEBNAD_BIN="/usr/bin/zebnad"
env ZEBNAD_USER="zebna"
env ZEBNAD_GROUP="zebna"
env ZEBNAD_PIDDIR="/var/run/zebnad"
# upstart can't handle variables constructed with other variables
env ZEBNAD_PIDFILE="/var/run/zebnad/zebnad.pid"
env ZEBNAD_CONFIGFILE="/etc/zebna/zebna.conf"
env ZEBNAD_DATADIR="/var/lib/zebnad"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # zebnad will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$ZEBNAD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run zebnad."
        echo "The setting must appear in $ZEBNAD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Zebna Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$ZEBNAD_PIDDIR"
    chmod 0755 "$ZEBNAD_PIDDIR"
    chown $ZEBNAD_USER:$ZEBNAD_GROUP "$ZEBNAD_PIDDIR"
    chown $ZEBNAD_USER:$ZEBNAD_GROUP "$ZEBNAD_CONFIGFILE"
    chmod 0660 "$ZEBNAD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$ZEBNAD_PIDFILE" \
    --chuid $ZEBNAD_USER:$ZEBNAD_GROUP \
    --exec "$ZEBNAD_BIN" \
    -- \
    -pid="$ZEBNAD_PIDFILE" \
    -conf="$ZEBNAD_CONFIGFILE" \
    -datadir="$ZEBNAD_DATADIR" \
    -disablewallet \
    -daemon

